# docker_emulator/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stderr_logfile=/app/logs/nginx.err.log
stdout_logfile=/app/logs/nginx.out.log

[program:xvfb]
command=Xvfb :99 -ac -screen 0 1280x720x24
autostart=true
autorestart=true
stderr_logfile=/app/logs/xvfb.err.log
stdout_logfile=/app/logs/xvfb.out.log

[program:gunserver]
command=node /app/gunServer.js
autostart=true
autorestart=true
stderr_logfile=/app/logs/gunserver.err.log
stdout_logfile=/app/logs/gunserver.out.log

[program:streaming]
command=/bin/bash /app/start_streaming.sh
autostart=true
autorestart=true
stderr_logfile=/app/logs/streaming.err.log
stdout_logfile=/app/logs/streaming.out.log
environment=PYTHONUNBUFFERED=1

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface